# config/file-service.service
[Unit]
Description=Personal File Upload Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/services
ExecStart=/home/YOUR_USERNAME/services/file_service_env/bin/python /home/YOUR_USERNAME/services/file_service.py
Restart=always
RestartSec=15
Environment=PATH=/home/YOUR_USERNAME/services/file_service_env/bin:/usr/local/bin:/usr/bin:/bin

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=file-service

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/mnt/storage

[Install]
WantedBy=multi-user.target
